# supervisord/mysqld.conf
[program:mysqld]
command=/usr/bin/pidproxy /var/mysqld/mysqld.pid /usr/bin/mysqld_safe --pid-file=/var/mysqld/mysqld.pid
autostart=true
autorestart=false
user=mysql
