# supervisord/nginx.conf
[program:nginx]
command=cat
autostart=true
autorestart=true
user=akaunting
