# supervisord/proxygen.conf
[program:proxygen]
command=/usr/bin/hhvm -m server -c /etc/hhvm/server.ini
autostart=true
autorestart=true
user=akaunting
